# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx -g "daemon off;"

[program:run]
command=sh /usr/share/nginx/html/run.sh

[program:runApache]
command=sh /usr/share/nginx/html/runApache.sh
